# server/systemd/jetlag-cas.service
[Unit]
Description=Jet Lag Hide and Seek PWA + CAS API
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/jetlag-hide-and-seek/server
Environment=CAS_DATA_DIR=/var/lib/jetlag-hide-and-seek
Environment=CAS_STATIC_DIR=/opt/jetlag-hide-and-seek/dist
Environment=CAS_STATIC_PREFIX=/JetLagHideAndSeek/
Environment=CAS_HOST=127.0.0.1
Environment=CAS_PORT=8787
Environment=NODE_ENV=production
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure

[Install]
WantedBy=multi-user.target
